# Set_With_Match_Confidence/build/devel/etc/catkin/profile.d/05.catkin-test-results.sh
# generated from catkin/cmake/env-hooks/05.catkin-test-results.sh.develspace.in

export CATKIN_TEST_RESULTS_DIR="/home/robotics/Desktop/OpenCV_Set/Set_With_Match_Confidence/build/test_results"
export ROS_TEST_RESULTS_DIR="$CATKIN_TEST_RESULTS_DIR"
